# lidarr/rootfs/etc/s6-overlay/s6-rc.d/lidarr/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Lidarr
# Runs Lidarr
# ==============================================================================

bashio::log.info 'Starting Lidarr...'

exec /opt/bin/Lidarr -nobrowser -data=/config
